# crontabs
# 每周1早上2:00
0 0 2 ? * MON cd /app && ./cli get >> /app/cron.log 2>&1
# 每周1早上2:30
0 30 2 ? * MON cd /app && ./cli rank --type=week >> /app/cron.log 2>&1
# 每个月1号早上3:00
0 0 3 1 * ? cd /app && ./cli rank --type=month >> /app/cron.log 2>&1
# 每年1月1日早上3:30
0 30 3 1 1 ? * cd /app && ./cli rank --type=year >> /app/cron.log 2>&1